# archiso/airootfs/etc/NetworkManager/dispatcher.d/01-no-send-hostname.sh
/etc/NetworkManager/dispatcher.d/no-wait.d/01-no-send-hostname.sh